# internal/files/tofuhut-workload@.timer
[Unit]
Description=Tofuhut workload %i timer

[Timer]
OnBootSec=10min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
